# src/translations/en.json
{
    "general": {
        "dashboard": "Dashboard",
        "system_cleaner": "System Cleaner",
        "gpu_info": "GPU Information",
        "gpu_information": "GPU Information",
        "system_repair": "System Repair",
        "dism_tool": "DISM Tool",
        "network_reset": "Network Reset",
        "disk_check": "Disk Check",
        "boot_repair": "Boot Repair",
        "virus_scan": "Virus Scan",
        "system_information": "System Information",
        "settings": "Settings",
        "save": "Save",
        "cancel": "Cancel",
        "back": "Back",
        "next": "Next",
        "finish": "Finish",
        "scan": "Scan",
        "clean": "Clean",
        "repair": "Repair",
        "check": "Check",
        "reset": "Reset"
    },
    "menu": {
        "file": "File",
        "exit": "Exit",
        "tools": "Tools",
        "system_info": "System Information",
        "help": "Help",
        "about": "About",
        "language": "Language"
    },
    "dashboard": {
        "system_status": "System Status",
        "cpu_usage": "CPU Usage",
        "memory_usage": "Memory Usage",
        "disk_usage": "Disk Usage",
        "system_temperature": "System Temperature",
        "quick_access": "Quick Access",
        "optimize_system": "Optimize System",
        "clean_junk": "Clean Junk Files",
        "scan_virus": "Scan for Viruses",
        "scan_system": "Scan System",
        "get_system_info": "Get System Info"
    },
    "system_cleaner": {
        "title": "System Cleaner",
        "description": "Clean unnecessary files to free up disk space",
        "scan_for": "Scan for:",
        "scan_options": "Scan Options",
        "temp_files": "Temporary files",
        "recycle_bin": "Recycle bin",
        "cache_files": "Cache files",
        "log_files": "Log files",
        "exclude_files": "Exclude files",
        "add_exclusion": "Add exclusion",
        "remove_exclusion": "Remove exclusion",
        "file_extensions": "File extensions",
        "add_extension": "Add extension",
        "remove_extension": "Remove extension",
        "scan_button": "Scan Now",
        "stop_button": "Stop",
        "clean_button": "Clean Selected",
        "clean_threats": "Clean Threats",
        "scanning": "Scanning...",
        "scan_stopped": "Scan stopped",
        "cleaning": "Cleaning...",
        "scan_complete": "Scan Complete",
        "clean_complete": "Clean Complete",
        "items_found": "Items Found",
        "space_to_free": "Space to Free",
        "scan_tab": "Scan",
        "results_tab": "Results",
        "settings_tab": "Settings",
        "ready_to_scan": "Ready to scan",
        "start_scan": "Start Scan",
        "clean_selected": "Clean Selected",
        "exclusions": "Exclusions",
        "add": "Add",
        "remove": "Remove",
        "confirm_action": "Confirm Action",
        "confirm_clean": "Are you sure you want to clean {count} detected files?"
    },
    "gpu_info": {
        "title": "GPU Information",
        "description": "View basic information about your GPU",
        "gpu_name": "GPU Name",
        "driver_version": "Driver Version",
        "memory_total": "Total Memory",
        "memory_used": "Memory Used",
        "temperature": "Temperature",
        "refresh": "Refresh Info"
    },
    "system_repair": {
        "title": "System Repair",
        "description": "Check and repair system issues",
        "scan_for": "Scan for:",
        "scan_options": "Scan Options",
        "registry_issues": "Registry issues",
        "system_files": "System files",
        "startup_items": "Startup items",
        "services": "Services",
        "scan_button": "Scan Now",
        "repair_button": "Repair Selected",
        "repair_selected": "Repair Selected",
        "scanning": "Scanning...",
        "repairing": "Repairing...",
        "scan_complete": "Scan Complete",
        "repair_complete": "Repair Complete",
        "issues_found": "Issues Found",
        "fixed_issues": "Fixed Issues",
        "scan_tab": "Scan",
        "results_tab": "Results",
        "settings_tab": "Settings",
        "ready_to_scan": "Ready to scan",
        "start_scan": "Start Scan",
        "clean_selected": "Clean Selected",
        "exclusions": "Exclusions",
        "add": "Add",
        "remove": "Remove"
    },
    "dism_tool": {
        "title": "DISM Tool",
        "description": "Deployment Image Servicing and Management Tool",
        "operations": "Operations:",
        "check_health": "Check Health",
        "scan_health": "Scan Health",
        "restore_health": "Restore Health",
        "cleanup_image": "Cleanup Image",
        "start_button": "Start Operation",
        "log_output": "Operation Log"
    },
    "network_reset": {
        "title": "Network Reset",
        "description": "Reset network settings to fix connection issues",
        "warning": "Warning: This will reset all network adapters and settings",
        "operations": "Operations:",
        "flush_dns": "Flush DNS Cache",
        "reset_winsock": "Reset Winsock",
        "reset_tcp_ip": "Reset TCP/IP Stack",
        "reset_firewall": "Reset Firewall Settings",
        "reset_button": "Reset Network",
        "log_output": "Operation Log"
    },
    "disk_check": {
        "title": "Disk Check",
        "description": "Check disk for errors and fix them",
        "select_drive": "Select Drive:",
        "check_types": "Check Types:",
        "file_system": "File System Errors",
        "bad_sectors": "Bad Sectors",
        "read_only": "Read-only Mode (does not fix issues)",
        "check_button": "Check Disk",
        "repair_button": "Repair Disk",
        "warning": "Warning: Disk check may require a system restart",
        "log_output": "Operation Log"
    },
    "boot_repair": {
        "title": "Boot Repair",
        "description": "Repair Windows startup issues. This tool can fix various startup issues, including损坏的MBR、BCD问题和丢失的启动文件。",
        "warning": "⚠️ WARNING: Boot repair operations may damage the system, please ensure you understand the purpose of each option and proceed with caution. These operations are simulated and will not affect the system.",
        "mbr": "Repair Master Boot Record (MBR)",
        "bcd": "Repair Boot Configuration Data (BCD)",
        "bootmgr": "Repair Boot Manager",
        "winload": "Repair Windows Loader",
        "stop": "Stop",
        "full": "Full Repair (All Options)",
        "repair_options": "Repair Options",
        "start_repair": "Start Repair",
        "complete": "Complete",
        "startup_description": "Manage automatically started programs. Disabling unnecessary startup items can improve startup speed and system performance.",
        "refresh_list": "Refresh List",
        "operations": "Operations:",
        "rebuild_bcd": "Rebuild BCD",
        "fix_mbr": "修复MBR",
        "fix_boot": "Repair Boot Sector",
        "repair_button": "Repair Boot",
        "log_output": "Operation Log",
        "startup_manager": "Startup Manager",
        "startup_items": "Startup Items",
        "enable": "Enable",
        "disable": "Disable",
        "delete": "Delete"
    },
    "virus_scan": {
        "title": "Virus Scan",
        "description": "Scan for viruses and malware",
        "scan_locations": "Scan Locations:",
        "system_drive": "System Drive",
        "all_drives": "All Drives",
        "specific_folder": "Specific Folder",
        "select_folder": "Select Folder",
        "scan_types": "Scan Types:",
        "quick_scan": "Quick Scan",
        "full_scan": "Full Scan",
        "custom_scan": "Custom Scan",
        "scan_button": "Start Scan",
        "stop_button": "Stop Scan",
        "scanning": "Scanning...",
        "scan_complete": "Scan Complete",
        "threats_found": "Threats Found",
        "clean_threats": "Clean Threats",
        "log_output": "Scan Log"
    },
    "settings": {
        "title": "Settings",
        "description": "Configure application settings",
        "general": "General",
        "language": "Language:",
        "theme": "Theme:",
        "theme_settings": "Theme Settings",
        "light": "Light",
        "dark": "Dark",
        "blue": "Blue Theme",
        "green": "Green Theme",
        "purple": "Purple Theme",
        "system": "System",
        "custom": "Custom",
        "background_color": "Background Color:",
        "text_color": "Text Color:",
        "accent_color": "Accent Color:",
        "choose_color": "Choose Color",
        "startup": "Start on boot",
        "updates": "Auto-check for updates",
        "advanced": "Advanced",
        "disk_space": "Minimum free disk space to maintain (GB):",
        "logs": "Keep logs for (days):",
        "exclusions": "Exclusions",
        "files_folders": "Excluded files and folders:",
        "add": "Add",
        "remove": "Remove",
        "save_button": "Save Settings",
        "reset_button": "Reset to Default"
    },
    "system_info": {
        "title": "System Information",
        "description": "View detailed information about your system",
        "hardware": "Hardware",
        "operating_system": "Operating System",
        "network": "Network",
        "applications": "Installed Applications",
        "cpu": "CPU",
        "memory": "Memory",
        "disk": "Disk",
        "graphics": "Graphics",
        "cpu_usage": "CPU Usage",
        "memory_usage": "Memory Usage",
        "disk_usage": "Disk Usage",
        "os_version": "OS Version",
        "environment_variables": "Environment Variables",
        "network_adapters": "Network Adapters",
        "ip_address": "IP Address",
        "mac_address": "MAC Address",
        "hostname": "Hostname",
        "refresh_button": "Refresh Information"
    }
} 